# systemd/hermes-sona.service
[Unit]
Description=Hermes SONA Daemon — Self-Optimizing Neural Architecture
Documentation=https://github.com/dlhiwig/hermes
After=network.target
Wants=network.target

[Service]
Type=simple
User=toba
WorkingDirectory=/home/toba/hermes
Environment=NODE_ENV=production
Environment=SONA_PORT=18805
EnvironmentFile=-/home/toba/hermes/.env
ExecStart=/home/linuxbrew/.linuxbrew/bin/node --experimental-vm-modules dist/brain/sona.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-sona

# Kill switch enforcement via resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
