# systemd/hermes-loop.service
[Unit]
Description=Hermes Loop — 8-Step Recursive Self-Learning Loop
Documentation=https://github.com/dlhiwig/hermes
After=network.target hermes-sona.service
Wants=hermes-sona.service

[Service]
Type=simple
User=toba
WorkingDirectory=/home/toba/hermes
Environment=NODE_ENV=production
EnvironmentFile=-/home/toba/hermes/.env
ExecStart=/home/linuxbrew/.linuxbrew/bin/node --experimental-vm-modules dist/main.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-loop

# Kill switch enforcement via resource limits
LimitNOFILE=65536
MemoryMax=4G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
